# aggregator/crontab
# cron は最小PATHで動くため /usr/local/bin (pythonの実体) を明示
PATH=/usr/local/bin:/usr/bin:/bin
# 30分ごとに収集 (FETCH_INTERVAL_MINに合わせて好きに変更)
# cron は環境変数を継承しないので、起動時にダンプした container.env を読み込む
*/30 * * * * cd /app && . /app/container.env && python main.py >> /app/data/aggregator.log 2>&1
# 必ず最後に改行が必要
